# Configer/supervisor_conf/dev_api.conf
[unix_http_server]
file=/var/supervisord/api_supervisor.sock   ; (the path to the socket file)

[inet_http_server]           ; inet (TCP) server disabled by default
port=*:9001                  ; (ip_address:port specifier, *:port for all iface)
username=user                ; (default is no username (open server))
password=123                 ; (default is no password (open server))

[supervisord]
logfile=/root/logs/supervisord/api_supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/var/supervisord/api_supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/supervisord/api_supervisor.sock ; use a unix:// URL  for a unix socket


;;;测试环境都是单例;;;
;;;AG Service层;;;

[program:AG_uwsgi]
command=uwsgi --chdir=/root/MrWind-Dispatcher/api_gateway/uwsgi --socket 127.0.0.1:4002 --wsgi-file /root/MrWind-Dispatcher/api_gateway/uwsgi/run_as_wsgi_app.py -p 1 --py-autoreload 1 --buffer-size=32768
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/root/log/ag/uwsgi.log

[program:AG_apps]
command=python /root/MrWind-Dispatcher/api_gateway/apps/run_server.py --port=4003
directory=/root/MrWind-Dispatcher/api_gateway/apps
redirect_stderr=true
stopsignal=QUIT
stdout_logfile=/root/log/ag/apps.log
loglevel=debug

;[program:AG_fence]
;command=python /root/MrWind-Dispatcher/api_gateway/fence/run_server.py --port=4005
;directory=/root/MrWind-Dispatcher/api_gateway/fence
;stopsignal=QUIT
;redirect_stderr=true
;loglevel=debug
;stdout_logfile=/root/log/ag/fence.log

;[program:AG_windchat]
;command=python /root/MrWind-Dispatcher/api_gateway/windchat/app.py --port=4006
;directory=/root/MrWind-Dispatcher/api_gateway/windchat
;stopsignal=QUIT
;redirect_stderr=true
;loglevel=debug
;stdout_logfile=/root/log/ag/windchat.log

;[program:AG_windchat_Q]
;command=python /root/MrWind-Dispatcher/api_gateway/windchat/windchat_queue.py
;directory=/root/MrWind-Dispatcher/api_gateway/windchat
;stopsignal=QUIT
;redirect_stderr=true
;loglevel=debug
;stdout_logfile=/root/log/ag/windchat/q_log.log

[program:AG_express]
command=python /root/MrWind-Dispatcher/api_gateway/express/run_server.py --port=4007 --logging=debug
directory=/root/MrWind-Dispatcher/api_gateway/express
stopsignal=QUIT
redirect_stderr=true
loglevel=debug
stdout_logfile=/root/log/ag/express/log.log



;;;BL Service层;;;

[program:BL_express]
command=python /root/MrWind-Dispatcher/business_logic/express/run_server.py --port=5002 --logging=debug
directory=/root/MrWind-Dispatcher/business_logic/express
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/root/log/bl/express.log

[program:SINGLETON_crontab]
command=python /root/MrWind-Dispatcher/business_logic/crontab/run_server.py
directory=/root/MrWind-Dispatcher/business_logic/crontab
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/root/log/bl/crontab.log

;[program:SINGLETON_coolie]
;command=python /root/MrWind-Dispatcher/business_logic/coolie/run_server.py --port=4008
;directory=/root/MrWind-Dispatcher/business_logic/coolie
;stopsignal=QUIT
;redirect_stderr=true
;stdout_logfile=/root/log/bl/coolie.log



;;;Data & Service层;;;

[program:DAS_number]
command=uwsgi --chdir=/root/MrWind-Dispatcher/data_and_service/number --socket 127.0.0.1:6001 --wsgi-file /root/MrWind-Dispatcher/data_and_service/number/run_as_wsgi_app.py -p 1 --async 100 --ugreen --enable-threads
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/root/log/das/number.log

[program:DAS_shop]
command=uwsgi --chdir=/root/MrWind-Dispatcher/data_and_service/shop --socket 127.0.0.1:6003 --wsgi-file /root/MrWind-Dispatcher/data_and_service/shop/run_as_wsgi_app.py -p 1 --py-autoreload 1 --buffer-size=32768
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/root/log/das/shop.log

[program:SINGLETON_apps]
command=python /root/MrWind-Dispatcher/data_and_service/apps/run_server.py --port=6004
directory=/root/MrWind-Dispatcher/data_and_service/apps
redirect_stderr=true
stopsignal=QUIT
stdout_logfile=/root/log/das/apps
loglevel=debug

;[program:DAS_fence]
;command=python /root/MrWind-Dispatcher/data_and_service/fence/run_server.py --port=6005
;directory=/root/MrWind-Dispatcher/data_and_service/fence
;stopsignal=QUIT
;redirect_stderr=true
;stdout_logfile=/root/log/das/fence.log

;[program:DAS_windchat]
;command=python /root/MrWind-Dispatcher/data_and_service/windchat/app.py --port=6006
;directory=/root/MrWind-Dispatcher/data_and_service/windchat
;stopsignal=QUIT
;redirect_stderr=true
;loglevel=debug
;stdout_logfile=/root/log/das/windchat.log

;[program:SINGLETON_gps]
;command=python /root/MrWind-Dispatcher/data_and_service/gps/run_server.py --port=6007
;directory=/root/MrWind-Dispatcher/data_and_service/gps
;redirect_stderr=true
;stopsignal=QUIT
;stdout_logfile=/root/log/das/gps.log
;loglevel=warn


[program:SINGLETON_push]
command=python -u /root/MrWind-Dispatcher/data_and_service/push/msg_server.py
process_name=%(process_num)s
numprocs=1
redirect_stderr=true
stopsignal=QUIT
stdout_logfile=/root/log/das/push%(process_num)s.log

[program:SINGLETON_sms]
command=python -u /root/MrWind-Dispatcher/data_and_service/sms/sms_server.py
redirect_stderr=true
stopsignal=QUIT
stdout_logfile=/root/log/das/sms.log

;[program:SINGLETON_default_mq]
;command=celery worker -A celery_worker -l info
;directory=/root/MrWind-Dispatcher/data_and_service/mq
;redirect_stderr=true
;stopsignal=QUIT
;stdout_logfile=/root/log/das/mq.log
