# Configer/supervisor_conf/prod_api.conf
[unix_http_server]
file=/root/mrwind/tmp/prod_api_supervisor.sock   ; (the path to the socket file)

[inet_http_server]           ; inet (TCP) server disabled by default
port=*:9001                  ; (ip_address:port specifier, *:port for all iface)
username=user                ; (default is no username (open server))
password=123                 ; (default is no password (open server))

[supervisord]
logfile=/root/mrwind/log/supervisord/prod_api_supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/root/mrwind/tmp/prod_api_supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///root/mrwind/tmp/prod_api_supervisor.sock ; use a unix:// URL  for a unix socket




;;;双例(2组配置);;;

[group:x000]
programs=0_AG_uwsgi,0_AG_express,0_BL_express,0_DAS_number,0_DAS_shop,0_AG_windchat,0_AG_windchat_Q,0_DAS_windchat

;;;AG Service层;;;

[program:0_AG_uwsgi]
command=uwsgi --chdir=/root/MrWind-Dispatcher/api_gateway/uwsgi --socket 127.0.0.1:4002 --wsgi-file run_as_wsgi_app.py -p 1 --async 2000 --ugreen --enable-threads
stopsignal=QUIT
stdout_logfile=NONE
stderr_logfile=NONE

[program:0_AG_express]
command=python /root/MrWind-Dispatcher/api_gateway/express/run_server.py --port=4003
directory=/root/MrWind-Dispatcher/api_gateway/express
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=NONE

; 4004 AG_apps占用;

;[program:0_AG_fence]
;command=python /root/MrWind-Dispatcher/api_gateway/fence/run_server.py --port=4005
;directory=/root/MrWind-Dispatcher/api_gateway/fence
;stopsignal=QUIT
;redirect_stderr=true
;stdout_logfile=/xvdc/log/ag/fence/0_access.log

[program:0_AG_windchat]
command=python /root/MrWind-Dispatcher/api_gateway/windchat/app.py --port=4006
directory=/root/MrWind-Dispatcher/api_gateway/windchat
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/xvdc/log/ag/windchat/0_access.log

;;; windchat Q 不需要端口 ;;;
[program:0_AG_windchat_Q]
command=python /root/MrWind-Dispatcher/api_gateway/windchat/windchat_queue.py
directory=/root/MrWind-Dispatcher/api_gateway/windchat
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/xvdc/log/ag/windchat_queue/0_access.log


;;;BL Service层;;;

[program:0_BL_express]
command=python /root/MrWind-Dispatcher/business_logic/express/run_server.py --port=5002
directory=/root/MrWind-Dispatcher/business_logic/express
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=NONE


;;;Data & Service层;;;

[program:0_DAS_number]
command=uwsgi --chdir=/root/MrWind-Dispatcher/data_and_service/number --socket 127.0.0.1:6001 --wsgi-file run_as_wsgi_app.py -p 1 --async 2000 --ugreen --enable-threads --disable-logging
stopsignal=QUIT
stdout_logfile=NONE
stderr_logfile=NONE

[program:0_DAS_shop]
command=uwsgi --chdir=/root/MrWind-Dispatcher/data_and_service/shop --socket 127.0.0.1:6003 --wsgi-file run_as_wsgi_app.py -p 1 --async 2000 --ugreen --enable-threads
stopsignal=QUIT
stdout_logfile=NONE
stderr_logfile=NONE

; 6004 DAS_apps占用;

;[program:0_DAS_fence]
;command=python /root/MrWind-Dispatcher/data_and_service/fence/run_server.py --port=6005
;directory=/root/MrWind-Dispatcher/data_and_service/fence
;stopsignal=QUIT
;redirect_stderr=true
;stdout_logfile=NONE

[program:0_DAS_windchat]
command=python /root/MrWind-Dispatcher/data_and_service/windchat/app.py --port=6006
directory=/root/MrWind-Dispatcher/data_and_service/windchat
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=NONE



;;; BELOW/ABOVE ARE BACKUP FOR ROLLBACK ;;;
;;; (1)加入新模块的时候要加入两个哦 ;;;
;;; (2)同时注意程序名,程序路径,端口号,log文件名要改成和该组一致哦 ;;;
;;; (3)别忘掉group里面也要添加哦 ;;;

[group:x100]
programs=1_AG_uwsgi,1_AG_express,1_BL_express,1_DAS_number,1_DAS_shop,1_AG_windchat,1_AG_windchat_Q,1_DAS_windchat

;;;AG Service层;;;

[program:1_AG_uwsgi]
command=uwsgi --chdir=/xvdb/MrWind-Dispatcher/api_gateway/uwsgi --socket 127.0.0.1:4102 --wsgi-file run_as_wsgi_app.py -p 1 --async 2000 --ugreen --enable-threads
stopsignal=QUIT
stdout_logfile=NONE
stderr_logfile=NONE


[program:1_AG_express]
command=python /xvdb/MrWind-Dispatcher/api_gateway/express/run_server.py --port=4103
directory=/xvdb/MrWind-Dispatcher/api_gateway/express
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=NONE

;[program:1_AG_fence]
;command=python /xvdb/MrWind-Dispatcher/api_gateway/fence/run_server.py --port=4105
;directory=/xvdb/MrWind-Dispatcher/api_gateway/fence
;stopsignal=QUIT
;redirect_stderr=true
;stdout_logfile=NONE

[program:1_AG_windchat]
command=python /xvdb/MrWind-Dispatcher/api_gateway/windchat/app.py --port=4106
directory=/xvdb/MrWind-Dispatcher/api_gateway/windchat
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=NONE

;;; windchat Q 不需要端口 ;;;
[program:1_AG_windchat_Q]
command=python /xvdb/MrWind-Dispatcher/api_gateway/windchat/windchat_queue.py
directory=/xvdb/MrWind-Dispatcher/api_gateway/windchat
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=NONE


;;;BL Service层;;;

[program:1_BL_express]
command=python /xvdb/MrWind-Dispatcher/business_logic/express/run_server.py --port=5102
directory=/xvdb/MrWind-Dispatcher/business_logic/express
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=NONE



;;;Data & Service层;;;

[program:1_DAS_number]
command=uwsgi --chdir=/xvdb/MrWind-Dispatcher/data_and_service/number --socket 127.0.0.1:6101 --wsgi-file run_as_wsgi_app.py -p 1 --async 2000 --ugreen --enable-threads --disable-logging
stopsignal=QUIT
stdout_logfile=NONE
stderr_logfile=NONE

[program:1_DAS_shop]
command=uwsgi --chdir=/xvdb/MrWind-Dispatcher/data_and_service/shop --socket 127.0.0.1:6103 --wsgi-file run_as_wsgi_app.py -p 1 --async 2000 --ugreen --enable-threads
stopsignal=QUIT
stdout_logfile=NONE
stderr_logfile=NONE

;[program:1_DAS_fence]
;command=python /xvdb/MrWind-Dispatcher/data_and_service/fence/run_server.py --port=6105
;directory=/xvdb/MrWind-Dispatcher/data_and_service/fence
;stopsignal=QUIT
;redirect_stderr=true
;stdout_logfile=NONE

[program:1_DAS_windchat]
command=python /xvdb/MrWind-Dispatcher/data_and_service/windchat/app.py --port=6106
directory=/xvdb/MrWind-Dispatcher/data_and_service/windchat
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/xvdc/log/das/windchat/1_access.log


;;;单例;;;

[program:SINGLETON_AG_apps]
command=python /xvdc/MrWind-Dispatcher/api_gateway/apps/run_server.py --port=4004
directory=/xvdc/MrWind-Dispatcher/api_gateway/apps
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=NONE
loglevel=warn

[program:SINGLETON_apps]
command=python /xvdc/MrWind-Dispatcher/data_and_service/apps/run_server.py --port=6004
directory=/xvdc/MrWind-Dispatcher/data_and_service/apps
redirect_stderr=true
stopsignal=QUIT
stdout_logfile=NONE
loglevel=warn

;[program:SINGLETON_gps]
;command=python /xvdc/MrWind-Dispatcher/data_and_service/gps/run_server.py --port=6007
;directory=/xvdc/MrWind-Dispatcher/data_and_service/gps
;redirect_stderr=true
;stopsignal=QUIT
;stdout_logfile=NONE
;loglevel=warn

[program:SINGLETON_crontab]
command=python /xvdc/MrWind-Dispatcher/business_logic/crontab/run_server.py
directory=/xvdc/MrWind-Dispatcher/business_logic/crontab
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/xvdc/log/bl/crontab/access.log

;[program:SINGLETON_coolie]
;command=python /xvdc/MrWind-Dispatcher/business_logic/coolie/run_server.py --port=6008
;directory=/xvdc/MrWind-Dispatcher/business_logic/coolie
;stopsignal=QUIT
;redirect_stderr=true
;stdout_logfile=NONE

;[program:SINGLETON_push]
;command=python -u /xvdc/MrWind-Dispatcher/data_and_service/push/msg_server.py
;process_name=%(process_num)s
;numprocs=1
;stopsignal=QUIT
;redirect_stderr=true
;stdout_logfile=/xvdc/log/das/push/access%(process_num)s.log

[program:SINGLETON_sms]
command=python -u /xvdc/MrWind-Dispatcher/data_and_service/sms/sms_server.py
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/xvdc/log/das/sms/access.log

;[program:SINGLETON_async_task]
;command=celery worker -A celery_worker -l info
;directory=/xvdc/MrWind-Dispatcher/data_and_service/mq
;stopsignal=QUIT
;redirect_stderr=true
;stdout_logfile=/xvdc/log/das/mq/access.log
